# deploy/docker/couchdb.conf
[program:couchdb]
autorestart=false
command=couchdb
environment=HOME=/usr/local/var/lib/couchdb
redirect_stderr=true
user=couchdb
priority=1
